# build-dev/openwrt/files/etc/uci-defaults/99-custom
#!/bin/sh

uci -q batch << EOI
set network.lan.proto='dhcp'
set luci.main.lang='zh_cn'
EOI
